# debian/62-Endless-enable-wayland.rules
# Eanble Wayland if there is the DRM
SUBSYSTEM=="graphics", KERNEL=="fb0", ATTR{name}!="simple", RUN+="/usr/lib/gdm3/gdm-disable-wayland FALSE"
