# debian/62-Endless-disable-wayland.rules
# disable Wayland
RUN+="/usr/lib/gdm3/gdm-disable-wayland"
